# src/8/dist/lib/systemd/system/appserver-php5-fpm.service

[Unit]
Description=appserver-php5-fpm php5-ts-fpm php-ts-fpm
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/appserver-php5-fpm.pid
ExecStart=/opt/appserver/sbin/php-fpm --php-ini /opt/appserver/etc/php-fpm-fcgi.ini --fpm-config /opt/appserver/etc/php-fpm.conf
ExecReload=/bin/kill -USR2 $MAINPID
PrivateTmp=true

[Install]
WantedBy=multi-user.target
